# config/pi_services/robo-path-api.service
[Unit]
Description=ROBO-Path FastAPI Storage Server
After=network.target

[Service]
Type=simple
User=rpi5
WorkingDirectory=%h/ROBO-Path_project
# '-' 접두사: .env 파일이 없어도 서비스 시작 실패하지 않음 (python-dotenv가 직접 로드)
EnvironmentFile=-%h/ROBO-Path_project/.env
ExecStart=%h/ROBO-Path_project/venv/bin/uvicorn src.infrastructure.storage.api:app --host 0.0.0.0 --port 8000
Restart=on-failure
RestartSec=5
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
