# config/pi_services/robo-path-api.service
[Unit]
Description=ROBO-Path FastAPI Storage Server
After=network.target

[Service]
Type=simple
User=rpi5
WorkingDirectory=/home/rpi5/ROBO-Path_project
EnvironmentFile=-/home/rpi5/ROBO-Path_project/.env
ExecStart=/home/rpi5/ROBO-Path_project/venv/bin/uvicorn src.infrastructure.storage.api:app --host 0.0.0.0 --port 8000
Restart=on-failure
RestartSec=5
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
